# supervisor/apache.conf
[program:apache2]
command=sigproxy /usr/sbin/apache2ctl -D FOREGROUND
priority=100
stdout_logfile=syslog
stderr_logfile=syslog
