# supervisor/apache.conf
[program:apache2]
command=pidproxy /var/run/apache2/apache2.pid /usr/sbin/apache2ctl -D FOREGROUND
priority=100
stdout_logfile=syslog
stderr_logfile=syslog
