# supervisor/cron.conf
[program:cron]
command=go-cron -file=/etc/go-cron/crontab
priority=100
directory=%(ENV_CRON_WORKDIR)s
