# supervisor/rsyslog.conf
[program:rsyslog]
command=rsyslogd -n
priority=10
